# scripts/chimcook.service
[Unit]
Description=Chimcook Voice Assistant
Documentation=https://github.com/SagitaKDX/Chimcook
# If docker isn't ready yet, retry a few times
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=%h/voice_assistant_v2

# Start the container stack
ExecStart=/usr/bin/docker compose up -d

# Stop cleanly on logout / shutdown
ExecStop=/usr/bin/docker compose stop

# If docker isn't ready yet, retry a few times
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=default.target
